# bin/zookeepr_test
#! /bin/sh

### BEGIN INIT INFO
# Provides:          zookeepr-test
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Zookeepr conference pylons webapp test server
# Description:       Zookeepr conference pylons webapp test server is launched by this script
#                    and runs as a script under paster
### END INIT INFO

SRV_DIR=/srv/http/test.followtheflow.org
PROJ_DIR=/home/contributors/repo/contributors
PID_FILE=$SRV_DIR/run/zookeepr_test.pid
LOG_FILE=$SRV_DIR/log/zookeepr_test.log
CONF_FILE=$PROJ_DIR/config.ini
USER=contributors
GROUP=lca2011

export PYTHONPATH=$PROJ_DIR/setup-dir

case "$1" in
  start)
    paster serve --user $USER --group $GROUP --daemon --pid-file=$PID_FILE --log-file=$LOG_FILE $CONF_FILE start
    ;;
  stop)
    paster serve --user $USER --group $GROUP --daemon --pid-file=$PID_FILE --log-file=$LOG_FILE  $CONF_FILE stop
    ;;
  restart)
    paster serve --user $USER --group $GROUP --daemon --pid-file=$PID_FILE --log-file=$LOG_FILE $CONF_FILE restart
    ;;
  *)
    echo $"Usage: $0 {start|stop|restart}"
    exit 1
esac 
